# gate/etc/cron.d/mail_update_aliases_shop
# Activity reports every 10 minutes everyday
*/1 * * * * wget -O - https://neo63.ru/service/alias_postfixadmin.php >/dev/null 2>/dev/null 
