# gate/etc/cron.d/mysqldump
# The first element of the path is a directory where the debian-sa1
# script is located
PATH=/usr/lib/sysstat:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin


# Additional run at 10:23 to rotate the statistics file
23 1 * * * root /home/scripts/mysql_dump.sh
